# devel/ansible/roles/common/files/synced-folder.service
[Unit]
Description=Wait for sshfs synced folder
After=network-online.target
Wants=network-online.target

[Service]
ExecStartPre=sh -c "while true; do [ -d /vagrant/.git ] && break; sleep 1; done"
ExecStart=sh -c "while true; do [ -d /vagrant/.git ] || exit 1; sleep 1; done"
Restart=always
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
