# mirrorlist-server/supervisor/mirrorlist-server.ini
[program:mirrorlist_server]
command=/usr/share/mirrormanager/mirrorlist-server/mirrorlist_server.py
priority=1
autostart=true
autorestart=true
startsecs=10
startretries=5
stopsignal=TERM
stopwaitsecs=10
user=mirrormanager
log_stderr=true
